# home/Library/LaunchAgents/homebrew.mxcl.rclone-mount.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>homebrew.mxcl.rclone-mount</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{ .rclone.binary }}</string>
        <string>mount</string>
        <string>{{ .rclone.remote }}:{{ .rclone.folder }}</string>
        <string>{{ .chezmoi.homeDir }}/{{ .rclone.mountPoint }}</string>
        <string>--daemon</string>
        <string>--vfs-cache-mode</string>
        <string>writes</string>
        <string>--vfs-cache-max-age</string>
        <string>1h</string>
        <string>--log-file</string>
        <string>{{ .rclone.logPath }}/rclone-mount.log</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
</dict>
</plist>
